# deploy/quadlet/vouchq-db.container
# Quadlet unit for the Postgres database (rootless systemd-managed Podman).
# Secrets come from an EnvironmentFile (see docs/self-hosting.md). The data dir
# is a named volume; the port is never published.
[Unit]
Description=Vouchq Postgres
Requires=vouchq-network.service
After=vouchq-network.service

[Container]
Image=docker.io/library/postgres:16-alpine
Network=vouchq.network
Volume=vouchq_pgdata.volume:/var/lib/postgresql/data:Z
EnvironmentFile=%h/.config/vouchq/db.env
HealthCmd=pg_isready -U ${POSTGRES_USER} -d ${POSTGRES_DB}
HealthInterval=10s
HealthRetries=10
NoNewPrivileges=true

[Service]
Restart=always

[Install]
WantedBy=default.target
